# infra/eureka/config/supervisord.conf
[supervisord]
nodaemon=true

[program:eureka]
command=/bin/bash -c "$CATALINA_HOME/bin/catalina.sh run"
redirect_stderr=true